# dist/br-thermal-autoadd.service
[Unit]
Description=scan for connected printers
After=br-thermal.service
Requires=br-thermal.service

[Service]
Type=oneshot
User=_br-thermal
Group=_br-thermal
SupplementaryGroups=plugdev
ExecStart=/usr/sbin/br-thermal autoadd
